# build/debian/howru-monitor_0.3.2/DEBIAN/postinst
#!/bin/sh

# Set file ownership
/bin/chown -R root:root /opt/howru
/bin/chown root:root /lib/systemd/system/howru-api.service
/bin/chown root:root /lib/systemd/system/howru-scheduler.service
/bin/chmod 655 /lib/systemd/system/howru-api.service
/bin/chmod 655 /lib/systemd/system/howru-scheduler.service
/bin/chmod 755 /opt/howru/gardener.py

# Enable service
systemctl daemon-reload
systemctl enable howru-scheduler
systemctl start howru-scheduler
systemctl enable howru-api.service
systemctl start howru-api.service

exit 0
