# demo/0.9/webapp/apps.conf
[supervisord]
nodaemon=true

[program:almond]
command=/opt/almond/almond
