# demo/0.8.3/custapp/apps.conf
[supervisord]
nodaemon=true

[program:almond]
command=/opt/almond/almond

[program:howru]
command=/opt/almond/howru

[program:custapp]
command=/custapp/app.py
