# demo/0.9/custapp/apps.conf
[supervisord]
nodaemon=true

[program:custapp]
command=/custapp/app.py
