# build/debian/howru-alerter_0.1.1/DEBIAN/postinst
#!/bin/sh

# Set file ownership
/bin/chown -R root:root /opt/howru
/bin/chown root:root /lib/systemd/system/howru-alerter.service
/bin/chmod 655 /lib/systemd/system/howru-alerter.service

# Enable service
systemctl daemon-reload
systemctl enable howru-alerter
systemctl start howru-alerter

exit 0
