# build/debian/almond-monitor-0.9.9/lib/systemd/system/howru.service
[Unit]
Description=HowRU api Service
After=multi-user.target network.target
Conflicts=getty@tty1.service

[Service]
User=root
Group=almond
WorkingDirectory=/opt/almond/www/api
ExecStart=/usr/bin/gunicorn --user almond --bind 0.0.0.0:80 --chdir=/opt/almond/www/api --log-file /var/log/almond/gunicorn.log --log-level info wsgi:app
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
StandardInput=tty-force

[Install]
WantedBy=multi-user.target
