# demo/0.8.3/kafka/kafka/supervisor/conf.d/almond.conf
[supervisord]
nodaemon=true

[program:almond]
command=/opt/almond/bin/almond

[program:howru]
command=/opt/almond/howru-api
