# build/debian/source/0.9.9.8/almond-monitor-0.9.9.8/lib/systemd/system/almond.service
[Unit]
Description=Almond Service
After=multi-user.target
Conflicts=getty@tty1.service

[Service]
Type=simple
User=almond
WorkingDirectory=/opt/almond/
ExecStart=/opt/almond/almond
StandardInput=tty-force

[Install]
WantedBy=multi-user.target
